# packaging/debian/etc/systemd/system/uptime-logger.service
[Unit]
Description=Log boot time to SQLite
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/uptime_logger

[Install]
WantedBy=multi-user.target
