# packaging/debian/DEBIAN/postinst
#!/bin/sh
set -e

echo "Creating data directory..."
mkdir -p /var/lib/uptime-logger

echo "Reloading systemd daemon..."
systemctl daemon-reload || true

echo "Enabling systemd services..."
systemctl enable uptime-logger.service || true
systemctl enable uptime-logger-shutdown.service || true

exit 0
