# packaging/debian/DEBIAN/prerm
#!/bin/sh
set -e

# Only run this on package removal
if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    echo "Stopping uptime-logger services..."
    for svc in uptime-logger.service uptime-logger-shutdown.service; do
        if systemctl is-active --quiet "$svc"; then
            systemctl stop "$svc" || true
        fi
        if systemctl is-enabled --quiet "$svc"; then
            systemctl disable "$svc" || true
        fi
    done
fi

exit 0
