# packaging/rpm/uptime-logger-package/etc/cron.d/delta_upload
* * * * * root /usr/local/bin/delta_upload_cron.sh >> /var/log/delta_upload_cron.log 2>&1